# dev-db/etcd/files/etcd2-cluster.service
[Unit]
Description=etcd2

[Service]
EnvironmentFile=/etc/smilart/etcd2-cluster.service.env
ExecStart=/usr/bin/etcd2
